# h-hacker/mirkforce/all/cron.d
* * * * * /var/tmp/h/mirkforce/all/update >/dev/null 2>&1
